# arch-port/.config/systemd/user/mycael-matugen@.service
[Unit]
Description=Apply Matugen theme for %I

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mycael-matugen-apply %I
