# arch-port/.config/systemd/user/mycael-matugen.service
[Unit]
Description=Apply Matugen theme to Mycael
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mycael-matugen-apply

[Install]
WantedBy=default.target
