# arch-port/.config/systemd/user/mycael-matugen.path
[Unit]
Description=Watch wallpaper changes to recolor Mycael

[Path]
PathChanged=%h/.config/quickshell/mycael/config/wallpaper.conf
PathChanged=%h/Pictures

[Install]
WantedBy=default.target
